# spark-worker-dirac/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[program:sshd]
autorestart=true
autostart=true
command=/usr/sbin/sshd -D
